# oui-ui-core/files/etc/uci-defaults/50_oui
#!/bin/sh
uci batch <<-EOF
	set uhttpd.main.ubus_prefix="/ubus"
	set uhttpd.main.index_page='oui.html'
	commit uhttpd
EOF
