# vwrt-ui-core/files/etc/uci-defaults/50_vuci
#!/bin/sh
uci batch <<-EOF
	set uhttpd.main.ubus_prefix="/ubus"
	set uhttpd.main.index_page='vuci.html'
	commit uhttpd
EOF
